# bin/gbd_diskspace.crontab
## Scripts para verificar os filesystems sgbd  - 2014-09-19 - Luis Marques
# Verificacao e envio de alertas
2,17,32,47 * * * * /postgresql/basedir/bin/gbd_diskspace.sh /postgresql/basedir/etc/gbd_diskspace.conf >> /postgresql/basedir/logs/gbd_diskspace_$(date -u +'\%Y-\%m').log 2>&1
# Reset do contador de alertas
1 9 * * 1 /postgresql/basedir/bin/gbd_diskspace_reset.sh >> /postgresql/basedir/logs/gbd_diskspace_$(date -u +'\%Y-\%m').log 2>&1
